# templates/docker.lvm/systemd/volume.mount.j2
[Unit]
Description={{ item.description | default(item.path)}}

[Mount]
What=/dev/{{ docker_lvm_vg }}/{{ item.lv }}
Where={{ item.path }}
Type=xfs

[Install]
WantedBy=local-fs.target