# templates/systemd-mountpoint.mount.j2
[Unit]
Description=Docker Root

[Mount]
What=/dev/{{ docker_lvm_vg }}/{{ docker_lvm_lv }}
Where=/var/lib/docker
Type=xfs

[Install]
WantedBy=local-fs.target